# cfg/cron.d/minecraft
# Shutdown minecraft (graceful) every 4th hour (shuts down only if no-one is playing)
0 */4 * * * minecraft /usr/bin/bash /home/minecraft/Fampack/cfg/scripts/shutdown.bash turnoff
